# auto-site-blocker/systemd/unblock-sites.timer
[Unit]
Description=Run site unblock at 2:00 AM daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
