# auto-site-blocker/systemd/unblock-sites.service
[Unit]
Description=Unblock YouTube and Reddit

[Service]
Type=oneshot
ExecStart=/usr/local/bin/unblock_sites.sh
